# docker/8.0/supervisord.d/getBTC.conf
[program:getBTC]
command=/usr/bin/php /var/www/html/artisan rates:getBTC
autostart=true
user=sail
stdout_logfile=/var/log/supervisor/default.log
