# docker/8.0/supervisord.d/schedule.conf
[program:schedule]
command=/usr/bin/php /var/www/html/artisan schedule:work
autostart=true
user=sail
stdout_logfile=/var/log/supervisor/default.log
